# player.service
[Unit]
Description=Raspberry Pi Random Video Player
After=network.target local-fs.target

[Service]
Type=simple
WorkingDirectory=/opt/rpi-random-player
Environment=USB_MOUNT_POINT=/mnt/usb
EnvironmentFile=-/etc/default/rpi-random-player
ExecStartPre=/bin/mkdir -p ${USB_MOUNT_POINT}
ExecStart=/usr/bin/python3 /opt/rpi-random-player/player.py
Restart=always
RestartSec=1
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
